# image.d/etc/services.d/odoo/run
#!/bin/execlineb -P

s6-setuidgid odoo
/srv/odoo/odoo-bin -c /srv/odoo/conf/odoo.conf
